# deploy/emberforge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.emberforge.backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/path/to/EmberForge/.venv/bin/emberforge</string>
    <string>serve</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/EmberForge</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>EMBERFORGE_ROOT</key>
    <string>/path/to/EmberForge</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/path/to/EmberForge/logs/emberforge.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/EmberForge/logs/emberforge.err.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>